# deployments/systemd/opengm-ca.service
# openGM-CA systemd服务配置
# 安装路径: /etc/systemd/system/opengm-ca.service

[Unit]
Description=openGM-CA Certificate Authority Service
Documentation=https://github.com/opengm-ca/opengm-ca
After=network-online.target openGauss.service
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/opengm-ca -config /etc/opengm-ca/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# 运行用户(必须预先创建)
User=opengm-ca
Group=opengm-ca

# 工作目录
WorkingDirectory=/var/lib/opengm-ca

# 环境变量
Environment="GOMAXPROCS=4"
Environment="GMCRYPTO_ENABLE_SM2=1"

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/opengm-ca /var/log/opengm-ca /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志输出到journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opengm-ca

[Install]
WantedBy=multi-user.target
